# debian/neutron-server.prerm
#!/bin/sh

set -e

case $1 in
        remove|purge)
                invoke-rc.d neutron-server stop
esac

#DEBHELPER#
