# debian/neutron-plugin-ml2-agent.upstart.in
description "Neutron ML2 plugin agent"
author "Thomas Goirand <zigo@debian.org>"

start on runlevel [2345]
stop on runlevel [016]

chdir /var/run

pre-start script
	mkdir -p /var/run/neutron
	chown neutron:root /var/run/neutron
end script

exec start-stop-daemon --start --chuid neutron:neutron --exec /usr/bin/neutron-ml2-agent -- --config-file=/etc/neutron/neutron.conf --config-file=/etc/neutron/plugins/ml2/ml2_conf.ini --log-file=/var/log/neutron/ml2-agent.log
